# services/s6-overlay/s6-rc.d/mosquitto-log/run
#!/bin/sh
sed 's/^/\x1b[34mmosquitto |\x1b[0m /' --unbuffered