# services/openrc/init.d/c8y-firmware-plugin
#!/sbin/openrc-run
description="Thin-edge device firmware management for Cumulocity"
command="/usr/bin/c8y-firmware-plugin"
: ${command_args=}
: ${command_user=tedge}
pidfile="/run/lock/${RC_SVCNAME}.lock"
TEDGE_LOGFILE="${TEDGE_LOGFILE:-/var/log/${RC_SVCNAME}.log}"
supervise_daemon_args=" --env TEDGE_RUN_LOCK_FILES=false"
error_log="${TEDGE_LOGFILE}"

depend() {
    after net
}
