# services/s6-overlay/s6-rc.d/tedge-configuration-plugin-log/run
#!/bin/sh
sed 's/^/\x1b[34mtedge-configuration-plugin |\x1b[0m /' --unbuffered