# services/supervisord/conf.d/c8y-configuration-plugin.conf
[program:c8y-configuration-plugin]
command=/usr/bin/c8y-configuration-plugin 
user=root
startsecs=5
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
