# services/supervisord/conf.d/tedge-agent.conf
[program:tedge-agent]
command=/bin/sh -c "tedge init ||:; /usr/bin/tedge-agent "
user=tedge
startsecs=5
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
