# services/supervisord/conf.d/tedge-mapper-c8y.conf
[program:tedge-mapper-c8y]
command=/bin/sh -c "tedge init ||:; /usr/bin/tedge-mapper c8y"
user=tedge
startsecs=5
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
