# services/supervisord/conf.d/tedge-log-plugin.conf
[program:tedge-log-plugin]
command=/usr/bin/tedge-log-plugin 
user=root
startsecs=5
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
