# services/supervisord/conf.d/tedge-configuration-plugin.conf
[program:tedge-configuration-plugin]
command=/usr/bin/tedge-configuration-plugin 
user=root
startsecs=5
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
