# services/s6-overlay/s6-rc.d/tedge-log-plugin-log/run
#!/bin/sh
sed 's/^/\x1b[34mtedge-log-plugin |\x1b[0m /' --unbuffered