# services/s6-overlay/s6-rc.d/tedge-mapper-az-log/run
#!/bin/sh
sed 's/^/\x1b[34mtedge-mapper-az |\x1b[0m /' --unbuffered