# services/s6-overlay/s6-rc.d/mosquitto/run
#!/command/with-contenv sh
exec 2>&1
exec /usr/sbin/mosquitto -c /etc/mosquitto/mosquitto.conf
